# etc/crontabs/root
*/10 * * * * ntpclient -h cn.pool.ntp.org -s -c 1
*/5 * * * * /etc/platform/bin/curlhttpgps.sh &
*/10 * * * * /etc/platform/bin/command.sh &
*/5 * * * * /etc/upgrade/rsync_task.sh &
#* * */1 * * /etc/platform/bin/fullUpdate.sh &
#* * */1 * * /etc/platform/bin/incrementalUpdate.sh &

